# debian/apache2-mpm-worker.prerm
#!/bin/sh -e

# Apache2 prerm
# Thom May <thom@debian.org>

case "$1" in
        upgrade|remove)
                /usr/sbin/invoke-rc.d apache2 stop
		sleep 2
        ;;
        deconfigure|failed-upgrade)
        ;;
        *)
                echo "prerm called with unknown argument \`$1'" >&2
                exit 1
        ;;
esac

exit 0
