# debian/apache2.2-common.preinst
#! /bin/sh -e

if [ "$1" = "upgrade" ] && dpkg --compare-versions "$2" lt 2.2; then
    if [ -x "/etc/init.d/apache2" ]; then
        if [ -x /usr/sbin/invoke-rc.d ]; then
            invoke-rc.d apache2 stop ||true
        else
            /etc/init.d/apache2 stop ||true
        fi
    fi
fi

#DEBHELPER#
